# templates/luistervink.cron
#birdnet
2-59/10 * * * * $USER sleep $(shuf -i 0-9 -n 1)m; /usr/local/bin/luistervink_task_check.sh >/dev/null 2>&1
